# package/mypkg/gateway-plc/files/gateway-plc.init
#!/bin/sh /etc/rc.common

START=97

SERVICE_USE_PID=1
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1
USE_PROCD=1

PLCGW=/usr/bin/gateway-plc



start_service() {

	procd_open_instance
	procd_set_param respawn 3600 5 0
	procd_set_param command $PLCGW
	procd_close_instance
}

stop_service() {
	service_stop $PLCGW
}


restart()
{  
	stop  
	start  
}  
